# samba.conffiles
/etc/cron.daily/samba
/etc/cron.weekly/samba
/etc/init.d/samba
